# geth.socket
[Unit]
Description=geth daemon (socket activation)

[Socket]
ListenStream=%t/geth.ipc


