# docker/supervisord.conf
[supervisord]
nodaemon=false

[program:geth]
priority=30
directory=/
command=geth --rpc
user=root
autostart=true
autorestart=true
startsecs=10
stopsignal=QUIT
stdout_logfile=/var/log/geth.log
stderr_logfile=/var/log/geth.err

[unix_http_server]
file=%(here)s/supervisor.sock

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface